# vpncheck.timer

#/etc/systemd/system/test.timer

[Unit]
Description=Timer for vpncheck

[Timer]
OnUnitActiveSec=1m
OnBootSec=1m

[Install]
WantedBy=timers.target
